# .config/systemd/user/autostart-wayland.target
[Unit]
Description=Current graphical user wayland session
Documentation=man:systemd.special(7)
RefuseManualStart=no
StopWhenUnneeded=no
