# .config/systemd/user/flameshot.service
[Unit]
Description=Flameshot Service

[Service]
Type=simple
ExecStart=/usr/bin/flameshot
Restart=on-failure

[Install]
WantedBy=autostart.target
