# .config/systemd/user/hypridle.service
[Unit]
Description=Hypridle Service

[Service]
Type=exec
ExecStart=/usr/bin/hypridle
Restart=on-failure

[Install]
WantedBy=autostart.target
