# .config/systemd/user/hyprlock.service
[Unit]
Description=Hyprlock Service

[Service]
Type=exec
ExecStartPre=-/usr/bin/playerctl --all-players pause
ExecStart=/usr/bin/hyprlock
