# .config/systemd/user/waybar.service
[Unit]
Description=Waybar Service

[Service]
Type=simple
ExecStart=/usr/bin/waybar
Restart=on-failure

[Install]
WantedBy=autostart-wayland.target
